# com.dougs.logfromwatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dougs.logfromwatch</string>

    <key>UserName</key>
    <string>dougs</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dougs/.local/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/dougs/PycharmProjects/LogFromWatch</string>
        <string>python</string>
        <string>/Users/dougs/PycharmProjects/LogFromWatch/main.py</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/logfromwatch.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/logfromwatch.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/dougs/PycharmProjects/LogFromWatch</string>
</dict>
</plist>
